# ansible/templates-injective/injectived.service.j2
[Unit]
Description=Injective Node ({{ node_name }})
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ ansible_user }}
Group={{ ansible_user }}
ExecStart=/usr/local/bin/injectived start --home {{ injective_home }}
Restart=on-failure
RestartSec=3
LimitNOFILE=65535
Environment="DAEMON_NAME=injectived"
Environment="DAEMON_HOME={{ injective_home }}"
Environment="DAEMON_ALLOW_DOWNLOAD_BINARIES=false"
Environment="DAEMON_RESTART_AFTER_UPGRADE=true"

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=injectived-{{ node_name }}

[Install]
WantedBy=multi-user.target

